# AlwaysAdd/replace_linux-preinstall_dir/etc/profile.d/linux-preinstall.sh
for DIR in \
    /opt/git/linux-preinstall/utilities-developer \
    /opt/git/linux-preinstall/utilities-server \
    /opt/git/linux-preinstall/utilities
do
    if [  -d $DIR ]; then
        PATH=$DIR:$PATH
    fi
done
